# supervisor_app.conf
[program:app-uwsgi]
command = uwsgi --ini /app/repairs_uwsgi.ini

[program:nginx-app]
command = /usr/sbin/nginx